# onvif-proxy.service
[Unit]
Description=ONVIF Proxy Service
Documentation=https://github.com/fpontieri/ONVIF-Proxy
After=network.target network-online.target
Wants=network-online.target
Requires=systemd-networkd.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/onvif-proxy
ExecStart=/usr/bin/python3 /opt/onvif-proxy/src/main_service.py --daemon --config /var/lib/onvif-proxy/config.xml
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=onvif-proxy

# Security settings
NoNewPrivileges=false
ProtectSystem=false
ProtectHome=false
ReadWritePaths=/opt/onvif-proxy /var/log /etc/systemd/network /var/lib/onvif-proxy /var/lib/dhcp /run /tmp
PrivateTmp=false
ProtectKernelTunables=false
ProtectKernelModules=false
ProtectControlGroups=false

# Network capabilities
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_DAC_OVERRIDE CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
